# deploy/llm-bot-dev@.service
[Unit]
Description=llm-bot dev (%i)
After=network-online.target
After=napcat@%i.service
Wants=network-online.target
Wants=napcat@%i.service
Wants=llm-bot-dev-webui@%i.service

[Service]
Type=simple
WorkingDirectory=%h/Workspace/nodejs/llm-onebot
Environment=NODE_ENV=development
Environment=CONFIG_INSTANCE=%i
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -lc '. %h/.nvm/nvm.sh && exec npm run dev'
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=default.target
